# supervisor/kafka.conf
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/bin/server.properties
autostart=true
autorestart=true